# data/autostart/com.github.dgsasha.Remembrance.desktop
[Desktop Entry]
Type=Application
Name=Remembrance
Exec=remembrance --gapplication-service
NoDisplay=true