# data/autostart/com.github.dgsasha.Remembrance.desktop.in
[Desktop Entry]
Type=Application
Name=@APP_NAME@
Exec=@PROJECT_NAME@ --gapplication-service
NoDisplay=true